# deploy/reviewer.service
# systemd user service for the local reviewer daemon.
#
# Install:
#   npm run build
#   mkdir -p ~/.config/systemd/user ~/.config/reviewer
#   # token (chmod 600):  echo 'GH_TOKEN=ghp_xxx' > ~/.config/reviewer/env
#   cp deploy/reviewer.service ~/.config/systemd/user/reviewer.service
#   # edit WorkingDirectory/REVIEWER_CONFIG below to your checkout path, then:
#   systemctl --user daemon-reload
#   systemctl --user enable --now reviewer
#   journalctl --user -u reviewer -f
#
# %h expands to the user's home directory.

[Unit]
Description=Local autonomous PR reviewer daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/projects/valkey-reviewer
Environment=REVIEWER_CONFIG=%h/projects/valkey-reviewer/reviewer.config.json
EnvironmentFile=-%h/.config/reviewer/env
ExecStart=/usr/bin/env node --enable-source-maps dist/daemon/src/cli.js daemon
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
