# deploy/revuto-guard.service
[Unit]
Description=revuto uptime guard
After=default.target

[Service]
Type=oneshot
ExecStart=%h/projects/revuto/scripts/revuto-guard.sh
